# build-rootfs/files/etc/init.d/podroid-bootstrap
#!/sbin/openrc-run
description="Podroid VM system bootstrap (cgroups, modules, devpts, sysctl)"

depend() {
    need localmount
    before podroid-network
}

start() {
    ebegin "Loading kernel modules..."
    echo "Loading kernel modules..." > /dev/console

    # Rootless podman needs / to be a shared mount so user namespaces can
    # propagate mount events. Without this, every rootless `podman run`
    # warns and may fail to mount /sys-as-fuse-overlayfs inside the userns.
    mount --make-rshared / 2>/dev/null

    # /dev sub-mounts (we replaced devfs OpenRC service with the initramfs's
    # devtmpfs + these explicit mounts — keeps the boot path tight)
    mkdir -p /dev/pts /dev/shm /dev/mqueue 2>/dev/null
    mountpoint -q /dev/pts    || mount -t devpts devpts /dev/pts \
        -o gid=5,mode=0620,ptmxmode=0666,noexec,nosuid
    mountpoint -q /dev/shm    || mount -t tmpfs tmpfs /dev/shm \
        -o noexec,nosuid,nodev,size=64m
    mountpoint -q /dev/mqueue || mount -t mqueue mqueue /dev/mqueue \
        -o noexec,nosuid,nodev
    [ -c /dev/ptmx ] || mknod /dev/ptmx c 5 2 2>/dev/null
    chmod 0666 /dev/ptmx /dev/pts/ptmx 2>/dev/null

    # Hostname (Alpine's hostname service is disabled in our runlevel)
    [ -r /etc/hostname ] && hostname -F /etc/hostname 2>/dev/null

    depmod -a 2>/dev/null
    for m in 9p 9pnet 9pnet_virtio; do modprobe "$m" 2>/dev/null; done
    mkdir -p /dev/net 2>/dev/null
    [ -c /dev/net/tun ] || mknod /dev/net/tun c 10 200 2>/dev/null
    [ -c /dev/fuse ]    || mknod /dev/fuse   c 10 229 2>/dev/null
    # Rootless containers need /dev/net/tun (pasta/slirp4netns) and /dev/fuse
    # (fuse-overlayfs) accessible from non-root users.
    chmod 0666 /dev/net/tun /dev/fuse 2>/dev/null

    # cgroup v2 subtree controllers
    mkdir -p /sys/fs/cgroup 2>/dev/null
    mountpoint -q /sys/fs/cgroup || mount -t cgroup2 cgroup2 /sys/fs/cgroup
    printf '+cpuset +cpu +io +memory +hugetlb +pids +rdma\n' \
        > /sys/fs/cgroup/cgroup.subtree_control 2>/dev/null

    # ZRAM swap (half RAM, lz4)
    if [ -b /dev/zram0 ]; then
        _mem_kb=$(awk '/^MemTotal:/{print $2}' /proc/meminfo)
        echo lz4 > /sys/block/zram0/comp_algorithm 2>/dev/null
        echo $((_mem_kb * 512)) > /sys/block/zram0/disksize 2>/dev/null
        mkswap /dev/zram0 >/dev/null 2>&1 && swapon -p 100 /dev/zram0 2>/dev/null
    fi

    # sysctl for container networking
    sysctl -qw net.ipv4.ip_forward=1 \
        net.ipv4.conf.all.forwarding=1 \
        net.bridge.bridge-nf-call-iptables=1 2>/dev/null

    eend 0
}
